# config/systemd/user/claude-process-issues@.timer
[Unit]
Description=process-issues for %i

[Timer]
OnBootSec=2min
OnUnitActiveSec=20min
RandomizedDelaySec=2min
Unit=claude-skill@process-issues:%i.service

[Install]
WantedBy=timers.target
